# remnawave-bot.service
[Unit]
Description=Remnawave Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/PromoRemnaBot
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/python3 /path/to/PromoRemnaBot/main.py
Restart=always
RestartSec=10

# Output to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remnawave-bot

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/path/to/PromoRemnaBot

[Install]
WantedBy=multi-user.target